# deploy/systemd/tko-jobs.timer
[Unit]
Description=TKOS RAG Platform - Periodic Jobs Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
